# firefly-iii/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Firefly III
# Runs the PHP-FPM daemon
# ==============================================================================

declare key

key=$(cat /data/firefly/appkey.txt)
export APP_KEY=${key}
export LOG_CHANNEL=stdout
export APP_ENV=local
export DB_CONNECTION=mysql
export DB_DATABASE
export DB_HOST
export DB_PASSWORD
export DB_PORT
export DB_USERNAME
export MAPBOX_API_KEY

if bashio::config.has_value 'mapbox_api_key';then
    MAPBOX_API_KEY=$(bashio::config "mapbox_api_key")
fi

if bashio::config.has_value 'remote_mysql_host';then
    DB_HOST=$(bashio::config "remote_mysql_host")
    DB_DATABASE=$(bashio::config "remote_mysql_database")
    DB_USERNAME=$(bashio::config "remote_mysql_username")
    DB_PASSWORD=$(bashio::config "remote_mysql_password")
    DB_PORT=$(bashio::config "remote_mysql_port")
else
    DB_HOST=$(bashio::services "mysql" "host")
    DB_DATABASE=firefly
    DB_USERNAME=$(bashio::services "mysql" "username")
    DB_PASSWORD=$(bashio::services "mysql" "password")
    DB_PORT=$(bashio::services "mysql" "port")
fi

if bashio::config.equals 'log_level' 'debug' \
    || bashio::config.equals 'log_level' 'trace' ;then
    export APP_DEBUG=true
    export APP_LOG_LEVEL=debug
fi
bashio::log.info "Installing/updating Database"
php /var/www/firefly/artisan firefly-iii:create-database > /dev/null 2>&1
php /var/www/firefly/artisan migrate --seed --no-interaction --force > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:decrypt-all > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:transaction-identifiers > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:migrate-to-groups > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:account-currencies > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:transfer-currencies > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:other-currencies > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:migrate-notes > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:migrate-attachments > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:bills-to-rules > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:bl-currency > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:cc-liabilities > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:back-to-journals > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:rename-account-meta > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:migrate-recurrence-meta > /dev/null 2>&1

php /var/www/firefly/artisan firefly-iii:fix-piggies > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:create-link-types > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:create-access-tokens > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:remove-bills > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:enable-currencies > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:fix-transfer-budgets > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:fix-uneven-amount > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:delete-zero-amount > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:delete-orphaned-transactions > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:delete-empty-journals > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:delete-empty-groups > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:fix-account-types > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:rename-meta-fields > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:fix-ob-currencies > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:fix-long-descriptions > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:fix-recurring-transactions > /dev/null 2>&1

php /var/www/firefly/artisan firefly-iii:report-empty-objects > /dev/null 2>&1
php /var/www/firefly/artisan firefly-iii:report-sum > /dev/null 2>&1

php /var/www/firefly/artisan passport:install > /dev/null 2>&1

php /var/www/firefly/artisan firefly-iii:set-latest-version --james-is-cool > /dev/null 2>&1
php /var/www/firefly/artisan cache:clear > /dev/null 2>&1
php /var/www/firefly/artisan config:cache > /dev/null 2>&1

bashio::log.info "Starting PHP-FPM..."

exec php-fpm7.4 --nodaemonize
