# firefly-iii/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Firefly III
# Runs the PHP-FPM daemon
# ==============================================================================

declare key

key=$(cat /data/firefly/appkey.txt)
export APP_KEY=${key}
export LOG_CHANNEL=stdout
export APP_ENV=local
if bashio::config.equals 'database' 'sqlite';then
    export DB_CONNECTION='sqlite'
fi
if bashio::config.equals 'database' 'mysql';then
    DB_HOST=$(bashio::config "mysql_host")
    DB_DATABASE=$(bashio::config "mysql_database")
    DB_USERNAME=$(bashio::config "mysql_user")
    DB_PASSWORD=$(bashio::config "mysql_password")
    export DB_HOST
    export DB_DATABASE
    export DB_USERNAME
    export DB_PASSWORD
fi
bashio::log.info "Installing/updating Database"
php /var/www/firefly/artisan package:discover
php /var/www/firefly/artisan migrate --seed
php /var/www/firefly/artisan firefly:decrypt-all
php /var/www/firefly/artisan firefly:upgrade-database
php /var/www/firefly/artisan passport:install
php /var/www/firefly/artisan cache:clear

bashio::log.info "Starting PHP-FPM..."

exec php-fpm7.2 --nodaemonize
